# openwrt/package/multimedia/hdpvrd/files/hdpvrd.init
#!/bin/sh /etc/rc.common
# Sample init script reading its
# configuration from /etc/config/package

START=50

start_service() {
	/bin/hdpvrd > /dev/null
}

